# rpi-monitor.service
[Unit]
Description=RPiMonitor — Raspberry Pi System Monitor
Documentation=https://github.com/coreconduit/rpi-monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/rpi-monitor
ExecStart=/usr/bin/python3 /opt/rpi-monitor/rpi_monitor.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment=PIMONITOR_HOST=0.0.0.0
Environment=PIMONITOR_PORT=8585
# Environment=PIMONITOR_TOKEN=your-secret-token
# Environment=PIMONITOR_SERVICES=ssh,nginx,docker,ollama,mosquitto
# Environment=PIMONITOR_REFRESH=2

# Hardening
ProtectSystem=strict
ReadWritePaths=/opt/rpi-monitor
ProtectHome=true
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
