# hub/rpi-monitor-hub.service
[Unit]
Description=RPiMonitor Hub — Lord of the Pi Monitors
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=coreconduit
WorkingDirectory=/opt/rpi-monitor/hub
ExecStart=/usr/bin/python3 /opt/rpi-monitor/hub/rpi_monitor_hub.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Optional env overrides
# Environment=PIHUB_PORT=8686
# Environment=PIHUB_TOKEN=your-secret
# Environment=PIHUB_POLL_INTERVAL=5
# Environment=PIHUB_DISCOVERY_PORT=8585

[Install]
WantedBy=multi-user.target
